# sys_files/usr/lib/systemd/system/horizon-user-setup.service
[Unit]
Description=Setup HorizonOS user aliases, completions and shell enhancements
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/horizon-user-setup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
