# sys_files/usr/lib/systemd/user/horizon-setup.service

[Unit]
Description=Setup HorizonOS user aliases, completions and shell enhancements
After=plasma-workspace.target
Wants=plasma-workspace.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/horizon-setup
RemainAfterExit=true

[Install]
WantedBy=default.target