# sys_files/usr/lib/systemd/user/horizon-setup-user.service

[Unit]
Description=Setup HorizonOS user aliases, completions and shell enhancements
After=plasma-workspace.target graphical-session.target
Wants=plasma-workspace.target graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/horizon-setup-user
RemainAfterExit=true

[Install]
WantedBy=default.target